# medications-mcp.service
[Unit]
Description=Unified Medication MCP Server
Documentation=https://github.com/zesty-genius128/Certus
After=network.target
Wants=network.target

[Service]
Type=simple
User=nodejs
Group=nodejs
WorkingDirectory=/opt/medication-mcp
ExecStart=/usr/bin/node hosted-mcp-server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/opt/medication-mcp/.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/medication-mcp/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=medication-mcp

[Install]
WantedBy=multi-user.target